# init.d/beep_sys
#! /bin/sh -e
### BEGIN INIT INFO
# Provides: beep_sys
# Required-Start: $syslog
# Required-Stop: $syslog
# Required-Reload: $syslog
# Default-Start: 2 3 4 5
# Default-Stop:	0 1 6
# Description: Emit a beep sound at basic system state changes
#
### END INIT INFO
#
#

set -e

case "$1" in
  start)
beep -n -f 600
;;
	stop)
beep -r 2 -l 10 -n -f 50 -r 2 -l 200
;;
	restart)
beep -f 100 -r 2 -l 10 -n -f 50 -r 2 -l 200 -n -f 600
;;
  *)
  echo "Use: /etc/init.d/beep_sys {start|stop|restart}"
  exit 1
  ;;
esac

exit 0
